# config/kanata/kanata.service
[Unit]
Description=Kanata keyboard remapper
Documentation=https://github.com/jtroo/kanata

[Service]
Type=simple
ExecStart=%h/.nix-profile/bin/kanata --cfg %h/.config/kanata/kanata.kbd
Restart=on-failure
RestartSec=3
Environment=PATH=/home/deck/.nix-profile/bin:/run/current-system/sw/bin:/usr/bin:/bin
Environment=DISPLAY=:0

[Install]
WantedBy=default.target
