# config/kanata/steamdeck/kanata.openrc
#!/sbin/openrc-run

# Kanata keyboard remapper OpenRC service script
name="kanata"
description="Kanata keyboard remapper"

command="/usr/bin/kanata"
command_args="--cfg /home/\${USER}/dotfiles/config/kanata/kanata.kbd"
command_user="\${USER}"
pidfile="/run/\${RC_SVCNAME}.pid"
command_background=true

depend() {
    need localmount
    use logger
}